# roles/developer/files/emacs/org.shtsh.emacs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>KeepAlive</key>
  <true/>
  <key>Label</key>
  <string>homebrew.mxcl.emacs</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/opt/emacs/bin/emacs</string>
    <string>--fg-daemon</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/emacs_sdout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/emacs_sderr.log</string>
</dict>
</plist>
